# deb/etc/systemd/system/gunshot.target
[Unit]
Description=Gunshot detector proxy and camera control functionality.

[Install]
WantedBy=multi-user.target